# flexget.service
#/etc/systemd/system/flexget.service
[Unit]
Description=Flexget

[Service]
Type=simple
ExecStart=/usr/bin/flexget execute
