# flexget.service
#/etc/systemd/system/flexget.service
[Unit]
Description=Flexget

[Service]
Type=simple
ExecStart=/bin/bash -c '/usr/local/bin/flexget execute'
